# personal-settings/usr/share/gvfs/mounts/network.mount
[Mount]
Type=network
Exec=/usr/lib/gvfsd-network
#AutoMount=true
AutoMount=true to AutoMount=false
